# crontab
#
#Add to crontab
#

*/1 *   * * *   root    /etc/cron.d/netflow-clean.sh
